# nas_ports/freenas/freenas-neighbord/pkg-plist
/usr/local/etc/rc.d/neighbord
/usr/local/sbin/neighbord
