# nas_ports/freenas/freenas-clid/pkg-plist
/usr/local/etc/rc.d/clid
/usr/local/sbin/clid
